# memory_layer.service
[Unit]
Description=SinoVec Memory Layer HTTP API
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=${SINOVEC_HOME}
EnvironmentFile=-/etc/default/sinovec
ExecStart=/usr/bin/python3 ${SINOVEC_HOME}/memory_layer.py serve --host 127.0.0.1 --port 18793
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
