# memory-sinovec.service
[Unit]
Description=Memory Sinovec HTTP API (pgvector + FastEmbed BAAI/bge-small-zh-v1.5)
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/SinoVec
EnvironmentFile=-/etc/default/sinovec
ExecStart=/usr/bin/python3 /opt/SinoVec/memory_sinovec.py serve --host 127.0.0.1 --port 18793
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
